# service/sync-on-server-ready.service
[Unit]
Description=Sync services when server name is reachable

[Service]
Type=notify
ExecStart=/usr/bin/sync-on-server-ready
